# etc/systemd/incidentrelay-telegram-worker.service
[Unit]
Description=IncidentRelay Telegram worker
After=network.target

[Service]
Type=simple
User=incidentrelay
Group=incidentrelay
WorkingDirectory=/var/www/incidentrelay

Environment=INCEDENTRELAY_CONFIG_FILE=/etc/incidentrelay/incidentrelay.conf
Environment=PYTHONUNBUFFERED=1

ExecStart=/usr/bin/python3 -m app.telegram_worker

Restart=always
RestartSec=5

KillSignal=SIGTERM
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
